# files/etc/service/pritunl-log/run
#!/bin/bash
tail --retry --follow /config/pritunl/log/pritunl.log
